# local-docker/hello-world/src/supervisor.conf
[supervisord]
nodaemon=true

[program:flask]
command=flask run -h 0.0.0.0 -p 5000

[program:apache2]
command=/usr/sbin/sshd -D -e